# beszel-healthcheck.timer
[Unit]
Description=Periodic Beszel health check

[Timer]
OnBootSec=5min
OnUnitActiveSec=90s

[Install]
WantedBy=timers.target
