# beszel-backup.service
[Unit]
Description=Beszel backup

[Service]
Type=oneshot
# Sync the entire data directory (PocketBase database and config) to the backup location
ExecStart=/usr/bin/rsync -a --delete %h/data/ /var/backups/beszel/
